# supervisord.conf
[supervisord]
logfile=./.relais/logs/supervisord.log
pidfile=./.relais/supervisord.pid
nodaemon=false
environment=DEBUGPY_ENABLED="%(ENV_DEBUGPY_ENABLED)s"

[unix_http_server]
file=./.relais/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix://./.relais/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; priority 1 — infrastructure
[program:courier]
command=redis-server %(here)s/config/redis.conf
priority=1
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=./.relais/logs/courier.log
redirect_stderr=true

; priority 8 — pure observers
[program:archiviste]
command=%(ENV_HOME)s/.local/bin/uv run python launcher.py archiviste/main.py
priority=8
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=./.relais/logs/archiviste.log
redirect_stderr=true
environment=PYTHONUNBUFFERED="1",PYTHONPATH="%(here)s",RELAIS_HOME="%(here)s/.relais",DEBUGPY_PORT="5681",DEBUGPY_WAIT="0"

; priority 10 — core bricks
[program:portail]
command=%(ENV_HOME)s/.local/bin/uv run python launcher.py portail/main.py
priority=10
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=./.relais/logs/portail.log
redirect_stderr=true
environment=PYTHONUNBUFFERED="1",PYTHONPATH="%(here)s",RELAIS_HOME="%(here)s/.relais",DEBUGPY_PORT="5679",DEBUGPY_WAIT="0"

[program:sentinelle]
command=%(ENV_HOME)s/.local/bin/uv run python launcher.py sentinelle/main.py
priority=10
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=./.relais/logs/sentinelle.log
redirect_stderr=true
environment=PYTHONUNBUFFERED="1",PYTHONPATH="%(here)s",RELAIS_HOME="%(here)s/.relais",DEBUGPY_PORT="5680",DEBUGPY_WAIT="0"

[program:atelier]
command=%(ENV_HOME)s/.local/bin/uv run python launcher.py atelier/main.py
priority=10
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=./.relais/logs/atelier.log
redirect_stderr=true
environment=PYTHONUNBUFFERED="1",PYTHONPATH="%(here)s",RELAIS_HOME="%(here)s/.relais",DEBUGPY_PORT="5678",DEBUGPY_WAIT="0"

[program:souvenir]
command=%(ENV_HOME)s/.local/bin/uv run python launcher.py souvenir/main.py
priority=10
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=./.relais/logs/souvenir.log
redirect_stderr=true
environment=PYTHONUNBUFFERED="1",PYTHONPATH="%(here)s",RELAIS_HOME="%(here)s/.relais",DEBUGPY_PORT="5682",DEBUGPY_WAIT="0"

[program:forgeron]
command=%(ENV_HOME)s/.local/bin/uv run python launcher.py forgeron/main.py
priority=10
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=./.relais/logs/forgeron.log
redirect_stderr=true
environment=PYTHONUNBUFFERED="1",PYTHONPATH="%(here)s",RELAIS_HOME="%(here)s/.relais",DEBUGPY_PORT="5685",DEBUGPY_WAIT="0"

[program:commandant]
command=%(ENV_HOME)s/.local/bin/uv run python launcher.py commandant/main.py
priority=10
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=./.relais/logs/commandant.log
redirect_stderr=true
environment=PYTHONUNBUFFERED="1",PYTHONPATH="%(here)s",RELAIS_HOME="%(here)s/.relais",DEBUGPY_PORT="5683",DEBUGPY_WAIT="0"

[program:horloger]
command=%(ENV_HOME)s/.local/bin/uv run python launcher.py horloger/main.py
priority=10
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=./.relais/logs/horloger.log
redirect_stderr=true
environment=PYTHONUNBUFFERED="1",PYTHONPATH="%(here)s",RELAIS_HOME="%(here)s/.relais",DEBUGPY_PORT="5686",DEBUGPY_WAIT="0"

; priority 20 — unified channel adapter supervisor
[program:aiguilleur]
command=%(ENV_HOME)s/.local/bin/uv run python launcher.py aiguilleur/main.py
priority=20
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=./.relais/logs/aiguilleur.log
redirect_stderr=true
environment=PYTHONUNBUFFERED="1",PYTHONPATH="%(here)s",RELAIS_HOME="%(here)s/.relais",DEBUGPY_PORT="5684",DEBUGPY_WAIT="0"

; priority 5 — WhatsApp Baileys gateway (external Node.js service)
[program:baileys-api]
command=%(ENV_HOME)s/.local/bin/uv run python scripts/run_baileys.py
directory=%(here)s
priority=5
autostart=false
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=./.relais/logs/baileys-api.log
redirect_stderr=true
environment=PYTHONUNBUFFERED="1",PYTHONPATH="%(here)s",RELAIS_HOME="%(here)s/.relais",REDIS_URL="redis://baileys:%(ENV_REDIS_PASS_BAILEYS)s@localhost:6379",PORT="3025",NODE_ENV="production"

; Groups
[group:core]
programs=portail,sentinelle,atelier,souvenir,commandant,archiviste,forgeron,horloger

[group:relays]
programs=aiguilleur

[group:infra]
programs=courier

[group:optional]
programs=baileys-api
