# supervisord.conf
[supervisord]
logfile=./.relais/logs/supervisord.log
pidfile=./.relais/supervisord.pid
nodaemon=false

[unix_http_server]
file=./.relais/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix://./.relais/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; priority 1 — infrastructure
[program:courier]
command=redis-server config/redis.conf
priority=1
autostart=true
autorestart=true
stdout_logfile=./.relais/logs/courier.log
redirect_stderr=true

; priority 8 — pure observers
[program:archiviste]
command=/Users/benjaminmarchand/.local/bin/uv run python archiviste/main.py
priority=8
autostart=true
autorestart=true
stdout_logfile=./.relais/logs/archiviste.log
redirect_stderr=true
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/Users/benjaminmarchand/IdeaProjects/relais",RELAIS_HOME="/Users/benjaminmarchand/IdeaProjects/relais/.relais"

; priority 10 — core bricks
[program:portail]
command=/Users/benjaminmarchand/.local/bin/uv run python portail/main.py
priority=10
autostart=true
autorestart=true
stdout_logfile=./.relais/logs/portail.log
redirect_stderr=true
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/Users/benjaminmarchand/IdeaProjects/relais",RELAIS_HOME="/Users/benjaminmarchand/IdeaProjects/relais/.relais"

[program:sentinelle]
command=/Users/benjaminmarchand/.local/bin/uv run python sentinelle/main.py
priority=10
autostart=true
autorestart=true
stdout_logfile=./.relais/logs/sentinelle.log
redirect_stderr=true
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/Users/benjaminmarchand/IdeaProjects/relais",RELAIS_HOME="/Users/benjaminmarchand/IdeaProjects/relais/.relais"

[program:atelier]
command=/Users/benjaminmarchand/.local/bin/uv run python -m debugpy --listen 5678 --wait-for-client atelier/main.py
priority=10
autostart=true
autorestart=true
stdout_logfile=./.relais/logs/atelier.log
redirect_stderr=true
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/Users/benjaminmarchand/IdeaProjects/relais",RELAIS_HOME="/Users/benjaminmarchand/IdeaProjects/relais/.relais"

[program:souvenir]
command=/Users/benjaminmarchand/.local/bin/uv run python souvenir/main.py
priority=10
autostart=true
autorestart=true
stdout_logfile=./.relais/logs/souvenir.log
redirect_stderr=true
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/Users/benjaminmarchand/IdeaProjects/relais",RELAIS_HOME="/Users/benjaminmarchand/IdeaProjects/relais/.relais"

; priority 20 — unified channel adapter supervisor
[program:aiguilleur]
command=/Users/benjaminmarchand/.local/bin/uv run python aiguilleur/main.py
priority=20
autostart=true
autorestart=true
stdout_logfile=./.relais/logs/aiguilleur.log
redirect_stderr=true
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/Users/benjaminmarchand/IdeaProjects/relais",RELAIS_HOME="/Users/benjaminmarchand/IdeaProjects/relais/.relais"

; Groups
[group:core]
programs=portail,sentinelle,atelier,souvenir,archiviste

[group:relays]
programs=aiguilleur

[group:infra]
programs=courier
